# scripts/systemd/sovereign-immich-external-restic.service
[Unit]
Description=Create a consistent Immich recovery copy on the external SSD
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sovereign-immich-external-restic backup
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
TimeoutStartSec=24h
