# scripts/systemd/sovereign-immich-protection@.service
[Unit]
Description=Sovereign Immich protection job (%i)
After=docker.service network-online.target
Wants=network-online.target
Requires=docker.service

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/sovereign-immich-protection %i
Nice=10
IOSchedulingClass=best-effort
IOSchedulingPriority=7
NoNewPrivileges=true
PrivateTmp=true
TimeoutStartSec=12h
