# scripts/systemd/sovereign-duckdns-update.timer
[Unit]
Description=Run Sovereign Homelab DuckDNS updater every 5 minutes

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min
AccuracySec=30s
Persistent=true
Unit=sovereign-duckdns-update.service

[Install]
WantedBy=timers.target
